# com.your_name_here.telegramcontrol.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.your_name_here.telegramcontrol</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/python@3.13/bin/python3.13</string>
        <string>/Users/your_name_here/scripts/telegram_control_bot.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/telegram_control_bot.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/telegram_control_bot_error.log</string>
</dict>
</plist>
